# private_Library/LaunchAgents/io.nettleton.brew-update.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.nettleton.brew-update</string>
  <!-- Daily gated brew updater (TESTING.md §7): safe-upgrade + reconcile +
       staleness cap. --upgrade-capped = the `mise run update:auto` variant;
       invoked directly so the job doesn't depend on mise activation.
       (Re)loaded by run_onchange_after_99-01_load-brew-update-agent when this
       file changes. pkg-installer casks can't sudo unattended — the updater
       detects that and posts a notification to run `mise run update`. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{ .chezmoi.homeDir }}/.local/share/chezmoi/tests/wrapper/daily_update.sh</string>
    <string>--upgrade-capped</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <!-- RunAtLoad: also fire at login and on plist reload — covers fire times
       missed while the machine was off or the agent not yet loaded. The
       updater is idempotent and cheap when there's nothing to do, and brew's
       own locking serializes any overlap with an in-flight apply. -->
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <!-- No ProcessType Background: it throttles the dozens of ruby/git
       processes brew spawns ~10-20x, stretching a ~3-minute run past an
       hour. Standard priority for a once-daily job is the better trade. -->
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/io.nettleton.brew-update.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/io.nettleton.brew-update.log</string>
</dict>
</plist>
